# cron.d/hebcron
# $Id$
# $URL$

SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=michael@radwin.org

# debug
# * * * * * nice touch /tmp/croncheck.hebcron.tmp

################## Hebcal ######################

HEBCAL_WEB=/home/hebcal/web/hebcal.com
HEBCAL_BIN=/home/hebcal/local/bin

# update static page at the beginning of each month
0 0 1 * * lynx -source 'http://www.hebcal.com/hebcal/?bogus=bogus' | perl -pi -e 's/\?bogus=bogus//g' > $HEBCAL_WEB/hebcal/default.new && mv -f $HEBCAL_WEB/hebcal/default.new $HEBCAL_WEB/hebcal/default.html

# update hebcal.com holiday.inc for front page at midnight
0 0 * * * nice $HEBCAL_WEB/bin/gen_current.pl

# update hebcal.com commentary at 11am on wednesday
11 11 * * 3 touch $HEBCAL_WEB/bin/calc_triennial.pl && sleep 2 && nice make -C $HEBCAL_WEB/bin > /dev/null

# shabbat: regular emails at 8:30pm pacific on Wednesdays
30 20 * * 3 nice $HEBCAL_BIN/shabbat_weekly.pl -all -verbose

# special Rosh Hashana 5771 - 7 Sep 2010 at 1:30am
#30 1 7 9 * nice $HEBCAL_BIN/shabbat_weekly.pl -all -verbose

# clear bounces at 3am on friday
3 3 * * 5 nice $HEBCAL_BIN/shabbat_deactivate.pl -count 5

# MySQL backups
30 0 * * * nice $HEBCAL_BIN/email_list_snapshot.sh

# SVN backup
13 0 * * * $HEBCAL_BIN/s3backup.sh

# clean cmd cache monthly
0 0 1 * * mv /tmp/hebcron-cache /tmp/hebcron-cache.$$ && rm -rf /tmp/hebcron-cache.$$

# don't run fetchmail as a daemon, run it from cron
*/6 * * * * /home/hebcal/local/bin/fetchmail
