# devops/etc/cron.d/hebcal
SHELL=/bin/sh
MAILTO="mradwin@hebcal.com"
HEBCAL_WEB=/var/www
HEBCAL_BIN=/home/hebcal/local/bin

# debug
#* * * * * hebcal touch /tmp/croncheck.hebcal.tmp

# monthly cache cleanup
0 0 1 * * www-data mv -f /tmp/hebcal-cache /tmp/hebcal-cache.$$ && rm -rf /tmp/hebcal-cache.$$

# update hebcal.com holiday.inc for front page at midnight
0 0 * * * hebcal $HEBCAL_WEB/bin/gen_current.pl

# update hebcal.com commentary at 11am on wednesday
11 11 * * 3 hebcal touch $HEBCAL_WEB/bin/calc_triennial.pl && sleep 2 && make -C $HEBCAL_WEB/bin > /dev/null 2>&1

# update Shabbat browse pages at 12:03am on Sundays
0 3 * * 0 hebcal $HEBCAL_WEB/bin/create_shabbat_browse.pl $HEBCAL_WEB/shabbat/browse

# MySQL backups
30 0 * * * hebcal $HEBCAL_BIN/email_list_snapshot.sh
